# external/opal-prd/opal-prd.service
[Unit]
Description=OPAL PRD daemon
ConditionVirtualization=false
ConditionPathExists=/sys/firmware/devicetree/base/ibm,opal/diagnostics
Documentation=man:opal-prd(8)

[Service]
EnvironmentFile=/etc/default/opal-prd
ExecStart=/usr/sbin/opal-prd $DAEMON_OPTS
Restart=always

[Install]
WantedBy=multi-user.target
